# app/docker/supervisor-app.conf
[program:app-uwsgi]
command = /root/.pyenv/shims/uwsgi --ini /home/nginx-config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
